# AllCommonScripts/RaspberryPiScripts/HomeScripts/.config/autostart/AVMediaStartup.desktop
[Desktop Entry]
Type=Application
Name=AVMediaStartup
Exec=/home/aprg/AVMediaStartup/startupAVMedia.sh
